# morze.service
[Unit]
Description=Morze tcp service

[Service]
ExecStart=/usr/bin/morze

[Install]
WantedBy=multi-user.target
